# systemd/douro.service
[Unit]
Description=Douro - Web hosting infrastructure analyzer with Prometheus export
Documentation=https://github.com/aboutmoi/Douro
After=network.target

[Service]
Type=simple
User=douro
Group=douro
WorkingDirectory=/opt/douro
Environment=PYTHONPATH=/opt/douro
ExecStart=/opt/douro/venv/bin/python -m douro.douro_exporter --config /opt/douro/config.production.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/douro
PrivateTmp=true

# Resource limits
MemoryMax=512M
TasksMax=10

[Install]
WantedBy=multi-user.target 